# systemd/hermes-overseer.service
[Unit]
Description=Hermes Overseer — VPS monitoring and response (distrobox)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
ExecStart=/usr/bin/distrobox enter hermes-overseer -- \
    %h/.local/share/hermes-overseer/.venv/bin/python -m overseer \
    --config %h/.config/hermes-overseer/overseer.yaml
EnvironmentFile=%h/.config/hermes-overseer/env
Restart=on-failure
RestartSec=30
WatchdogSec=600

[Install]
WantedBy=default.target
